# core/systemd/user/no1r_reset_budget.service
[Unit]
Description=Reset no1r daily budget (simulation)

[Service]
Type=oneshot
ExecStart=/home/sntrblck/.openclaw/workspace/scripts/reset_daily_budget.sh
User=sntrblck
