# core/systemd/user/no1r_heal_monitor.service
[Unit]
Description=no1r heal monitor service

[Service]
Type=oneshot
ExecStart=/home/sntrblck/.openclaw/workspace/scripts/monitor_no1r_heal.sh
User=sntrblck
