# core/systemd/user/paper_pilot_report.service
[Unit]
Description=Summarize paper trading pilot

[Service]
Type=oneshot
ExecStart=/home/sntrblck/.openclaw/workspace/scripts/summarize_paper_pilot.sh
User=sntrblck
